# pi/pi-sensor.conf
# pi sensor reporter
#
# pi sensor

description	"pi sensor"

start on net-device-up IFACE=eth0

umask 022

console log

script
    /usr/local/pi-sensor/daemon
end script
